# server.conf
[program:server]
command=e:/work/hchq/server.py
autorestart=true